# sv/service/pipewire-pulse/run
#!/bin/sh
exec 2>&1
SVDIR="$XDG_CONFIG_HOME/sv/enabled" sv start pipewire

# FIXME: chrt fails here due to insufficent permissions
chrt -r -p 80 $$
exec dbus-launch pipewire -c pipewire-pulse.conf
